# plotter.cron
*/15 * * * * apache cd /var/www; git pull; crontab -r; crontab plotter.cron; python plotter.py pfgplots.sqlite3
