# plotter.cron
*/15 * * * * apache cd /var/www; git pull; python plotter.py pfgplots.sqlite3
