# plotter.cron
*/15 * * * * apache plotter.py pfgplots.sqlite3
