# plotter.cron
*/15 * * * * apache python plotter.py pfgplots.sqlite3
